# src/program-redis-server.conf
[program:redis-server]
command = redis-server /srv/docker-registry/redis.conf
